# new-setup/conf/startupsound.service

[Unit]
Description=Startup sound

[Service]
ExecStart=/usr/bin/startup-sound.sh

[Install]
WantedBy=multi-user.target

